# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/tmp/supervisord.pid

[program:go-backend]
command=/app/prophet_bot
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=3
redirect_stderr=true
stdout_logfile=/app/logs/go-backend.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=HOME="/app",GIN_MODE="release"

[program:node-agent]
command=node agent/server.js
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=5
redirect_stderr=true
stdout_logfile=/app/logs/node-agent.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=HOME="/app",NODE_ENV="production"
