# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/tmp/supervisord.pid

; Node agent manages Go backend lifecycle internally via startGoBackend()
; Do NOT run go-backend as a separate supervisord program — port 4534 conflict

[program:node-agent]
command=node agent/server.js
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=10
redirect_stderr=true
stdout_logfile=/app/logs/node-agent.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=HOME="/app",NODE_ENV="production",GIN_MODE="release"
